# system_files/usr/lib/systemd/system/bash-config-setup.service
[Unit]
Description=Tweak default .bashrc
After=local-fs.target
ConditionPathExists=!/var/lib/bash-config-setup.done

[Service]
Type=oneshot
ExecStart=/usr/bin/cat >> %h/.bashrc << EOF
# ghostty
if [ -n "${GHOSTTY_RESOURCES_DIR}" ]; then
    builtin source "${GHOSTTY_RESOURCES_DIR}/shell-integration/bash/ghostty.bash"
fi
if [[ "$TERM" == "xterm-ghostty" ]]; then
    # Check if the terminfo actually exists; if not, fallback to xterm-256color
    if ! infocmp xterm-ghostty >/dev/null 2>&1; then
        export TERM=xterm-256color
    fi
fi

# brew
if [ -d "/home/linuxbrew/.linuxbrew/bin" ]; then
    eval "$(/home/linuxbrew/.linuxbrew/bin/brew shellenv)"
elif [ -d "$HOME/.linuxbrew/bin" ]; then
    # Fallback for single-user installation in home directory
    eval "$($HOME/.linuxbrew/bin/brew shellenv)"
fi

# mise
if command -v mise &> /dev/null; then
    eval "$(mise activate bash)"
fi

# Enable Flatpak exports
if [ -d "/var/lib/flatpak/exports/bin" ]; then
    export PATH="$PATH:/var/lib/flatpak/exports/bin"
fi
EOF
ExecStartPost=/usr/bin/touch /var/lib/bash-config-setup.done
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
